# Asn2/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# Virtual framebuffer providing an X display for Swing
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1024x768x24 -nolisten tcp
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Minimal window manager so the JFrame gets borders/title
[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# VNC server attached to Xvfb display
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport 5900 -nopw -quiet
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# noVNC (WebSocket <-> VNC bridge) + static web client on port 6080
[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# The Swing game itself. Restarts automatically when the user closes the window.
# gosu drops from root to uid 1000 (app user) so the JVM runs unprivileged,
# which satisfies the Kyverno "run-as-nonroot" policy for the main workload process.
[program:playgame]
command=/bin/sh -c "exec gosu app java -cp /app/classes PlayGame ${BOARD_SIZE} ${EMPTY_POSITIONS} ${DEPTH}"
directory=/app
environment=DISPLAY=":99",HOME="/app"
autorestart=true
startretries=1000
startsecs=0
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
